# lib/bindings/stitchingbox/etc/systemd/system/videostitch.service.d/videostitchenv.conf
[Service]
Environment="LIBVA_DRIVERS_PATH=/opt/intel/mediasdk/lib64"
Environment="LIBVA_DRIVER_NAME=iHD"
Environment="MFX_HOME=/opt/intel/mediasdk"
Environment="LD_LIBRARY_PATH=/opt/videostitch/lib/release:${LD_LIBRARY_PATH}:/opt/intel/mediasdk/lib64"